# monitor_test.timer
[Unit]
Description=Monitor every minute

[Timer]
OnCalendar=*:0/1
Unit=test-monitor.service

[Install]
WantedBy=timers.target
